# RUN/calibrator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Enabled</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string></string>
	</dict>
	<key>RunAtLoad</key>
	<false/>
	<key>KeepAlive</key>
	<true/>

	<key>ProgramArguments</key>
	<array>
	  <string>/Users/bst/Dropbox/Pulsefield/src/modules/Calibrator/calibrator</string>
	  <string>-d2</string>
	  <string>-D</string>
	  <string>/Users/bst/Desktop/PF_LOGS/calibrator_dbg.log</string>
	</array>
	<key>Label</key>
	<string>com.tc.calibrator</string>
	<key>WorkingDirectory</key>
	<string>/Users/bst/Dropbox/Pulsefield/src/modules/Calibrator</string>
	<key>StandardOutPath</key>
	<string>/Users/bst/Desktop/PF_LOGS/calibration_stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/bst/Desktop/PF_LOGS/calibration_stderr.log</string>
</dict>
</plist>
