# RUN/gpufluid.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Enabled</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string></string>
	</dict>
	<key>RunAtLoad</key>
	<false/>
	<key>KeepAlive</key>
	<true/>

	<key>Label</key>
	<string>com.tc.gpuflids</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Users/bst/Dropbox/Pulsefield/src/Video/Tracker/OF/GPUFluids/bin/GPUFluidsDebug.app/Contents/MacOS/GPUFluidsDebug</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/bst/Dropbox/Pulsefield/src/RUN</string>
	<key>StandardOutPath</key>
	<string>/Users/bst/Desktop/PF_LOGS/gpufluids_stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/bst/Desktop/PF_LOGS/gpufluids_stderr.log</string>
</dict>
</plist>
