# RUN/tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Enabled</key>
	<false/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string></string>
	</dict>
	<key>RunAtLoad</key>
	<false/>
	<key>KeepAlive</key>
	<true/>

	<key>Label</key>
	<string>com.tc.tracker</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/bin/java</string>
		<string>-Djava.library.path=/Users/bst/Dropbox/Pulsefield/src/Video/Tracker/bin/runTracker_lib</string>
		<string>-jar</string>
		<string>/Users/bst/Dropbox/Pulsefield/src/Video/Tracker/bin/runTracker.jar</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/bst/Dropbox/Pulsefield/src/Video/Tracker/bin</string>
	<key>StandardOutPath</key>
	<string>/Users/bst/Dropbox/Pulsefield/src/RUN/tracker_stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/bst/Dropbox/Pulsefield/src/RUN/tracker_stderr.log</string>
</dict>
</plist>
